# quadlets/filebrowser-quantum.container
[Unit]
Description=FileBrowser Quantum (Rootless Podman)

[Container]
ContainerName=filebrowser-quantum
Image=docker.io/gtstef/filebrowser:latest
AutoUpdate=registry

# --- Volumes ---
Volume=filebrowser-quantum_data:/home/filebrowser/data
Volume=filebrowser-quantum_files:/srv

# --- Time / Timezone ---
Volume=/etc/localtime:/etc/localtime:ro
Volume=/etc/timezone:/etc/timezone:ro
Environment=TZ=Europe/Berlin

# --- Admin Password ---
Environment=FILEBROWSER_ADMIN_PASSWORD="admin"

# --- Network ---
# Make sure the network exists before starting
Network=proxy-filebrowser-quantum

# --- Ports ---
# Internal container port 80 is used by FileBrowser
# No need to publish host ports when using Nginx Proxy Manager

[Service]
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
