# quadlets/opencloud/dev-wopi.container
[Unit]
Description=WOPI Server (Rootless Podman)
After=opencloud.service
Requires=opencloud.service

[Container]
ContainerName=wopi
Image=docker.io/opencloudeu/opencloud-rolling:latest

# Network
Network=proxy-opencloud

# Volumes
Volume=${OC_CONFIG_DIR:-/home/srvpi/opencloud-config}:/etc/opencloud:Z

# Environment
Environment=COLLABORATION_GRPC_ADDR=0.0.0.0:9301
Environment=COLLABORATION_HTTP_ADDR=0.0.0.0:9300
Environment=MICRO_REGISTRY=nats-js-kv
Environment=MICRO_REGISTRY_ADDRESS=opencloud:9233

# URLs angepasst für Nginx Proxy Manager + DuckDNS Domains
Environment=COLLABORATION_WOPI_SRC=https://wopi.example.org/
Environment=COLLABORATION_APP_NAME=CollaboraOnline
Environment=COLLABORATION_APP_PRODUCT=Collabora
Environment=COLLABORATION_APP_ADDR=https://collabora.example.org/
Environment=COLLABORATION_APP_ICON=https://collabora.example.org/favicon.ico
Environment=COLLABORATION_APP_INSECURE=true
Environment=COLLABORATION_CS3API_DATAGATEWAY_INSECURE=true
Environment=COLLABORATION_LOG_LEVEL=info
Environment=OC_URL=https://opencloud.example.org/

# Entrypoint / Exec
Entrypoint=/bin/sh
Exec=-c 'opencloud collaboration server'

[Service]
Restart=always
RestartSec=30s

[Install]
WantedBy=default.target
