# quadlets/filebrowser-quantum/onlyoffice.container
[Unit]
Description=Onlyoffice (Rootless Podman)
# Systemd service description

[Container]
# Container name
ContainerName=onlyoffice

# Container image
Image=docker.io/onlyoffice/documentserver:latest

# Load the timezone from Host
# Mount host timezone files for correct log timestamps
Volume=/etc/localtime:/etc/localtime:ro
Volume=/etc/timezone:/etc/timezone:ro

# Configure OnlyOffice for document editing:
Environment=JWT_ENABLED=true
Environment=JWT_SECRET="SUPER_SECRET_PASSWORD"

# --- Network ---
# Container network configuration
# IMPORTANT: The network must exist before starting this service!
# The container should be attached to the same network as Nginx Proxy Manager,
# similar to how Filebrowser-aquantum is connected.
# Example:
# podman network create npm-filebrowser-network
# podman network connect npm-filebrowser-network nginx-proxy-manager
Network=npm-filebrowser-network

# --- Ports ---
# Internal container port 80 is used by FileBrowser
# Container name (filebrowser-quantum) is selected in Nginx Proxy Manager as Forward Hostname
# Forward Port in Nginx Proxy Manager should be 80
# No need to publish host ports when using Nginx Proxy Manager
# PublishPort=80   # NOT required

[Service]
# Restart the container automatically if it crashes
Restart=always
RestartSec=10

[Install]
# Start the service at boot
WantedBy=default.target
